# RunAuto/crontabs/root
# test文件定时测试。可在模块目录看到每两分钟有执行一次的 test.log（用来测试定时启动是否正常，把它注释掉就关闭了）
# 所有的定时规则即改即刻生效，无需重启

*/2 * * * * /data/adb/modules/autofrp/RunAuto/sh/test.sh

# frpc定时启动。每50分钟运行一次(可自行修改，修改即可生效，无需重启)本人红米10x网络掉线frpc也跟着死了，所以必须结合定时启动使用。还有系统的一些其他因素。这里看个人需求，更多定时规则自己百度

*/50 * * * * /data/adb/modules/autofrp/RunAuto/sh/frpc.sh

# 热点状态 (默认关闭)由于手机不可控因素，偶尔出现热点自动关闭现象（关闭了还打不开，需要飞行模式），故每三分钟检测一次热点是否开启，如果关闭，则开启飞行模式，搭配xposededge使用(当然也可以通过命令)可实现热点常开(天玑处理器专用)

# */3 * * * * /data/adb/modules/autofrp/RunAuto/sh/hotspot_status.sh